# snapserver/rootfs/etc/services.d/mpd/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Runs mpd
# ==============================================================================
bashio::log.info 'Starting mpd...'

# Run mpd
exec mpd --no-daemon
